# crontab
@reboot set-gtk-theme
*/15 * * * * set-gtk-theme
